# silabs-multiprotocol-ember/rootfs/etc/s6-overlay/s6-rc.d/socat-zigbeed-tcp/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start socat TCP client for Zigbeed daemon
# ==============================================================================

ezsp_listen_port="$(bashio::addon.port 9999)"

# If user port is not set, listen on local interface only
if ! bashio::var.has_value "${ezsp_listen_port}"; then
    ezsp_listen_port="9999"
fi

bashio::log.info "Starting socat TCP client for Zigbeed daemon..."
exec  "/usr/bin/socat" -d \
      -L /var/run/socat-zigbeed-tcp \
      pty,link=/tmp/ttyZigbeeNCP,rawer,unlink-close=0 \
      tcp-l:${ezsp_listen_port},reuseaddr,fork