# silabs-multiprotocol-ember/rootfs/etc/s6-overlay/s6-rc.d/cpcd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Co-Processor Communication Daemon (CPCd)
# ==============================================================================

bashio::log.info "Starting cpcd..."
exec "/usr/bin/stdbuf" -o0 /usr/local/bin/cpcd
